# dist/etc/init.d/harold-bot
#!/bin/bash
### BEGIN INIT INFO
# Provides: harold-bot
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 1
# Short-Description: harold-bot for Concur Blue Slack group
### END INIT INFO

# This will prevent initlog from swallowing up a pass-phrase prompt if
# mod_ssl needs a pass-phrase from the user.
INITLOG_ARGS=""

USER=harold-bot
EXEC=/opt/harold-bot/bin/www
PIDFILE=/var/run/harold-bot.pid

. /lib/lsb/init-functions
test -x $EXEC || exit 0

# ========================================
# Options
# ----------------------------------------

case "$1" in
  start)
    log_daemon_msg "Starting  harold-bot"

    # Ensure there is some sort of config file present
    if [ ! -f /etc/harold-bot/config.sh ]; then
        if [ ! -f /etc/harold-bot/config.sh.default ]; then
            echo "ERROR: config file missing"
            exit 1
        fi
    fi

    # Load the default config
    if [ -f /etc/harold-bot/config.sh.default ]; then
        source /etc/harold-bot/config.sh.default
    fi

    # Overwrite the default values with any custom
    if [ -f /etc/harold-bot/config.sh ]; then
        source /etc/harold-bot/config.sh
    fi

    if start-stop-daemon --start --quiet --oknodo --make-pidfile --background --chuid $USER --pidfile $PIDFILE \
                    --startas /bin/bash -- -c "exec $EXEC | logger -i -t HaroldBot"; then

        log_end_msg 0
    else
        log_end_msg 1
    fi
    ;;

  stop)
    log_daemon_msg "Stopping  harold-bot"

    if start-stop-daemon --stop --quiet --oknodo --pidfile $PIDFILE; then
        log_end_msg 0
    else
        log_end_msg 1
    fi
    ;;

  status)
        status_of_proc -p $PIDFILE $EXEC harold-bot && exit 0 || exit $?
        ;;
  *)
        log_action_msg "Usage: /etc/init.d harold-bot {start|stop|status}"
        exit 1
esac

exit 0

